# scripts/crontab
# Watchman Agent Cron Schedule
# Runs Nemesis 1.0 daily

# Run Nemesis 1.0 every day at 8 AM UTC
0 8 * * * cd /app && PYTHONPATH=/app/scripts python3 scripts/nemesis/run_nemesis.py >> /data/logs/nemesis.log 2>&1

# Clean up old reports (keep last 30 days)
0 10 1 * * find /data/reports -name "*.json" -mtime +30 -delete
